# deploy/ehc-helpdesk.service
[Unit]
Description=EHC Helpdesk FastAPI Server
After=network.target ehc-vllm.service
Wants=ehc-vllm.service

[Service]
Type=simple
User=phungkien
WorkingDirectory=/home/phungkien/EHC_HELPDESK/ehc-helpdesk
Environment="PATH=/home/phungkien/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/home/phungkien/.local/bin/uvicorn api.routes:app --host 0.0.0.0 --port 8080
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
